# os/live-build/config/includes.chroot/etc/systemd/system/odysseus-service.service
[Unit]
Description=Odysseus AI agent (FastAPI, loopback :7000)
After=foulfox-prepare.service network-online.target
Wants=network-online.target
Requires=foulfox-prepare.service

[Service]
Type=simple
User=foulfox
Group=foulfox
EnvironmentFile=/etc/foulfox/foulfox.env
# Run from the active app release (a symlink the patcher flips for live updates);
# defaults to the baked /opt/foulfox/app. See foulfox-first-run + foulfox-patcher.
WorkingDirectory=/var/lib/foulfox/app-current/artifacts/odysseus-service
ExecStart=/usr/bin/env bash /var/lib/foulfox/app-current/artifacts/odysseus-service/start.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
